# deploy/ubuntu/cerium-backend.service
[Unit]
Description=Cerium Sales Tools API (Node.js/Express)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/backend/.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cerium-backend

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
